# platforms/linux_server/deb/postinst
#!/bin/bash
chown root:root /etc/systemd/system/jimbernetworkisolation.service
systemctl enable jimbernetworkisolation
systemctl start jimbernetworkisolation